# firstboot.service
[Unit]
Description=firstboot
ConditionPathExists=/firstboot.sh

[Service]
Type=oneshot
ExecStart=/firstboot.sh

[Install]
WantedBy=multi-user.target
